# packages/api/config/systemd/kit_power_down.service
[Unit]
Description=power down screen
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target

[Service]
ExecStart=/usr/bin/node -e 'require("/home/pi/intrakit/dist/utils/system.util").powerDown()'
Type=oneshot

[Install]
WantedBy=halt.target poweroff.target