# packages/api/config/systemd/kit_start_app.service
[Unit]
Description=start app
After=network.target

[Service]
EnvironmentFile=/home/pi/intrakit/production.env
ExecStart=/usr/bin/node /home/pi/intrakit/dist/main

[Install]
WantedBy=multi-user.target