# packages/api/systemd/i2c_detect.service
[Unit]
Description=Runs i2c detect script to track devices
Wants=i2c_detect.timer

[Service]
Type=oneshot
ExecStart=/home/ty/intrakit/scripts/i2cdetect.sh

[Install]
WantedBy=multi-user.target